# driver-pulsar/deploy/templates/bookkeeper.service
[Unit]
Description=Bookkeeper
After=network.target

[Service]
Environment="BOOKIE_MEM=-Xms4g -Xmx4g -XX:MaxDirectMemorySize=4g"

ExecStart=/opt/pulsar/bin/pulsar bookie
WorkingDirectory=/opt/pulsar
RestartSec=15s
Restart=on-failure
Type=simple
LimitNOFILE=300000

[Install]
WantedBy=multi-user.target
